# setup/messages/fr/crontab
# Journalier, supprimer les vieux fichier lock
50 1 * * * $INSTALL_DIR/bin/mj_trigger -t lock
# Horaire, effectue certaines tâches
20 * * * * $INSTALL_DIR/bin/mj_trigger -t hourly

