# setup/messages/en/crontab
# Daily, remove old lock files.
50 1 * * * $INSTALL_DIR/bin/mj_trigger -t lock
# Hourly, perform various tasks.
20 * * * * $INSTALL_DIR/bin/mj_trigger -t hourly

